# test-agent.service
[Unit]
Description=ThemeSmith Test Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ThemeSmith
ExecStart=/usr/bin/node /root/ThemeSmith/test_agent.js
Restart=always
RestartSec=10
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
